# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
loglevel=info
pidfile=/var/run/supervisord.pid

; ─── nginx (serves frontend + reverse-proxies API) ─────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log
priority=10

; ─── .NET Identity & Payment Service (port 5001) ──────────────────────────
[program:dotnet]
command=/app/net/net-service
directory=/app/net
environment=ASPNETCORE_URLS="http://localhost:5001",
            ASPNETCORE_ENVIRONMENT="Production",
            JWT__KEY="%(ENV_JWT_KEY)s"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/dotnet.err.log
stdout_logfile=/var/log/supervisor/dotnet.out.log
priority=20

; ─── Node.js Realtime & LMS Service (port 3001) ───────────────────────────
[program:nodejs]
command=node /app/njs/dist/index.js
directory=/app/njs
environment=PORT="3001",
            NODE_ENV="production",
            AGORA_APP_ID="%(ENV_AGORA_APP_ID)s",
            AGORA_APP_CREDENTIAL="%(ENV_AGORA_APP_CREDENTIAL)s",
            NET_SERVICE_URL="http://localhost:5001"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nodejs.err.log
stdout_logfile=/var/log/supervisor/nodejs.out.log
priority=20
